# packages/old/athena-system-installation/etc/profile.d/01-expand_cowspace.sh
#!/bin/bash

sudo mount -o remount,size=2G /run/archiso/cowspace
