# packages/archive/old/athena-system-installation/usr/lib/systemd/system/athena-graphical-target.service
[Unit]
Description=Athena Graphical Target
Before=systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/athena-graphical-target

[Install]
WantedBy=multi-user.target
