# packages/archive/old/athena-system-installation/etc/xdg/autostart/athena-welcome.desktop
[Desktop Entry]
Name=Athena Welcome
GenericName=Athena Welcome
X-GNOME-FullName=Athena Welcome
Comment=Welcome Application for Athena
Exec=/usr/local/bin/athena-welcome
Icon=athena-hello-logo
Terminal=false
Type=Application
Categories=GTK;Settings;Security;X-GNOME-Settings-Panel;X-GNOME-SystemSettings;X-Unity-Settings-Panel;X-XFCE-SettingsDialog;X-XFCE-SystemSettings;
Keywords=welcome;tool;
StartupNotify=true
