# archiso/airootfs/etc/profile.d/10_strap-once.sh
#!/bin/sh

if [ "$(whoami)" != "liveuser" ]; then
   /usr/local/bin/strap.sh   
fi
