# archiso/airootfs/etc/systemd/system/plymouth-quit.service.d/super_long_splash.conf
[Unit]
Description=Longer Plymouth Boot Screen time

[Service]
ExecStartPre=/usr/bin/sleep 5
